# system/restic-backup.service
[Unit]
Description=restic backup (personal-apps + stack) to Azure Blob
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic-backup.sh
Nice=10
IOSchedulingClass=idle
